# dotfiles/.config/systemd/user/emacs.service
[Unit]
Description=Emacs

[Service]
Type=simple
ExecStart=/run/current-system/sw/bin/emacs --fg-daemon
ExecStop=/run/current-system/sw/bin/emacsclient --eval "(kill-emacs)"
Environment=SSH_AUTH_SOCK=%t/keyring/ssh
Restart=always

[Install]
WantedBy=default.target
